# docker/supervisord-web.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700
chown=www-data:www-data

[supervisord]
nodaemon=true
logfile=/var/www/html/storage/logs/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/var/www/html/storage/logs

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/www/html/storage/logs/nginx.err.log
stdout_logfile=/var/www/html/storage/logs/nginx.out.log
priority=10

[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/www/html/storage/logs/php-fpm.err.log
stdout_logfile=/var/www/html/storage/logs/php-fpm.out.log
priority=20

[program:post-startup]
command=/bin/sh /var/www/html/docker/post-startup.sh
autostart=true
autorestart=false
startretries=1
startsecs=0
exitcodes=0
stderr_logfile=/var/www/html/storage/logs/post-startup.err.log
stdout_logfile=/var/www/html/storage/logs/post-startup.out.log
priority=99
